# deploy/systemd/wiki-arcana-monitor.service
[Unit]
Description=Wiki Arcana external health monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
DynamicUser=true
StateDirectory=wiki-arcana-monitor
RuntimeDirectory=wiki-arcana-monitor
LoadCredential=vault-role-id:/etc/wiki-arcana-monitor/vault-role-id
LoadCredential=vault-secret-id:/etc/wiki-arcana-monitor/vault-secret-id
Environment=VAULT_ROLE_ID_FILE=%d/vault-role-id
Environment=VAULT_SECRET_ID_FILE=%d/vault-secret-id
Environment=OPSBOT_VAULT_PATH=arcanada/data/prod/wiki-arcana/opsbot
ExecStart=/usr/bin/bash /srv/apps/wiki-arcana/current/scripts/monitor-wiki.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
